# onichip-backend.service
[Unit]
Description=OnichipGPS Backend API Server
Documentation=https://github.com/onichip/backend
After=network.target mongod.service
Wants=network.target
Requires=mongod.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/onichip/backend
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=MONGODB_URI=mongodb://localhost:27017/onichip_gps
Environment=JWT_SECRET=onichip_super_secret_key_production_2024
Environment=CORS_ORIGIN=http://localhost:80
ExecStart=/usr/bin/node servidor.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
KillMode=process
TimeoutSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=onichip-backend

# Límites de recursos para AWS t2.micro/t3.micro
LimitNOFILE=65536
LimitNPROC=4096
MemoryHigh=400M
MemoryMax=500M

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/onichip/backend
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
